# transcoding-api.service
[Unit]
Description=Local Transcoding API Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/api-storage.arkturian.com
Environment="PYTHONUNBUFFERED=1"
ExecStart=/var/www/api-storage.arkturian.com/venv/bin/python -m uvicorn transcoding_api_server:app --host 0.0.0.0 --port 8082
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transcoding-api

[Install]
WantedBy=multi-user.target
